# systemd/claude-pty.service.tmpl
[Unit]
Description=Claude PTY broker for {{project_name}} (E2 console)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{os_user}}
Group={{os_user}}
SupplementaryGroups=autodev-mcp
WorkingDirectory=/opt/{{slug}}

# RuntimeDirectory creates /run/claude-pty under owner=User mode 0755; the
# broker chgrp's the socket itself to autodev-mcp inside cli.js so the HTTP
# bridge (running as user autodev-mcp) can connect.
RuntimeDirectory=claude-pty
RuntimeDirectoryMode=0755

Environment=PTY_SOCKET_PATH=/run/claude-pty/{{slug}}.sock
Environment=PTY_SOCKET_GROUP=autodev-mcp
Environment=PTY_CWD=/opt/{{slug}}
Environment=PTY_COMMAND=/home/{{os_user}}/.local/bin/claude
Environment=PTY_COLS=120
Environment=PTY_ROWS=32

# The broker binary is shipped by the autodev-mcp service, not the project.
# Projects share one broker image; only the unit + env vars are per-slug.
ExecStart=/usr/bin/node /opt/autodev-mcp/dist/pty/cli.js
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=10

# Hardening — broker just spawns claude in a PTY, no syscall surface beyond
# fork+exec+socket I/O.
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/{{slug}} /run/claude-pty
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
# claude needs HTTPS to api.anthropic.com plus the unix socket for clients.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
