# infrastructure/deploy/praxis-hygiene.timer
[Unit]
Description=Daily Praxis disk hygiene
Documentation=https://github.com/g-chappell/praxis/blob/main/docs/runbooks/deploy-orchestrator.md

[Timer]
# Daily at 04:00 UTC, after the 03:30 Postgres backup. Persistent=true catches up
# if the VPS was off when the timer was due.
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=20m
Unit=praxis-hygiene.service

[Install]
WantedBy=timers.target
