# infrastructure/deploy/praxis-web.service
[Unit]
Description=Praxis web (Next.js) container
Documentation=https://github.com/g-chappell/praxis/blob/main/docs/runbooks/deploy-web.md
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
Restart=on-failure
RestartSec=5
TimeoutStopSec=20

# Pull the latest image before each (re)start so `systemctl restart`
# after a deploy picks up the new tag.
ExecStartPre=-/usr/bin/docker rm -f praxis-web
ExecStartPre=/usr/bin/docker pull ghcr.io/g-chappell/praxis-web:latest

# Host port :3002 (Colonize uses :3000, Pirate-Battle uses :3001 — see
# ADR-0004 port-allocation note). Container's internal port stays :3000
# (Next.js default; bound by the Dockerfile's EXPOSE + HOSTNAME env).
# Bound to 127.0.0.1 so only Caddy reaches it.
#
# --env-file passes /etc/praxis/praxis.env into the container at runtime
# (DATABASE_URL for @praxis/db, BETTER_AUTH_SECRET/URL for Better Auth,
# RESEND_API_KEY/FROM for the mailer). See docs/runbooks/deploy-web.md.
# --network praxis-net puts the container on the same bridge as
# praxis-db, so DATABASE_URL can use `praxis-db` as the hostname.
# Docker's built-in DNS resolves container names within a custom
# network. See docs/runbooks/deploy-postgres.md.
ExecStart=/usr/bin/docker run --rm --name praxis-web \
  --network praxis-net \
  --env-file /etc/praxis/praxis.env \
  -p 127.0.0.1:3002:3000 \
  ghcr.io/g-chappell/praxis-web:latest

ExecStop=/usr/bin/docker stop praxis-web

[Install]
WantedBy=multi-user.target
